# templates/default/redis.upstart.erb
# THIS FILE HAS BEEN GENERATED BY CHEF
# ANY MANUAL MODIFICATIONS WILL BE OVERWRITTEN

description "redis"

start on runlevel [23]
stop on shutdown

exec sudo -u redis /usr/local/bin/redis-server <%= node[:redis][:config] %>

respawn
